# scripts/systemd/gw-blocked-remediation-watch.service
[Unit]
Description=Groupware Kanban blocked remediation inotify watcher
After=default.target

[Service]
Type=simple
WorkingDirectory=/home/wrhrgw/gw
Environment=MAX_DISPATCH=1
Environment=EVENT_COOLDOWN_SECONDS=10
ExecStart=/usr/bin/python3 /home/wrhrgw/gw/scripts/gw-blocked-remediation-inotify-watch.py --board groupware
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
